# vm-exercise/systemd/gex-worker.service
# EXAMPLE unit — adapt paths/user. Runs the stateful worker under gunicorn.
[Unit]
Description=GEX VM sample - stateful worker (trading engine analog)
After=network.target gex-gateway.service
Wants=gex-gateway.service

[Service]
Type=simple
User=gexvm
WorkingDirectory=/opt/gex-sample-vm/service
Environment=DB_PATH=/var/lib/gexvm/data/service.db
Environment=GATEWAY_URL=http://127.0.0.1:8071
Environment=GATEWAY_TOKEN=REPLACE_FROM_SECRET_MANAGER
Environment=PORT=8070
ExecStart=/opt/gex-sample-vm/.venv/bin/gunicorn --bind 127.0.0.1:8070 --workers 1 worker:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
